# minecraft.service
[Unit]
After=network.service

[Service]
ExecStart=/opt/minecraft/minecraft.sh
User=minecraft

[Install]
WantedBy=default.target
