# minecraft-update-ip.service
[Unit]
After=network.service

[Service]
ExecStart=/opt/minecraft/minecraft-update-ip.sh

[Install]
WantedBy=default.target
